# crontab
* * * * * root /root/stream-producer.sh > /var/log/cron.log 2>&1

